# bin/iris
#!/bin/bash
#
# Init file for IRIS server daemon
#
# chkconfig: 2345 97 02
# description: IRIS server daemon
#
# processname: iris

# source function library
. /etc/rc.d/init.d/functions

RETVAL=0
SERVICE=$(basename $0)
PROG=java_service
DIR=/usr/share/java/iris-server
USER=tms
OPTIONS=$"java -server -Xmx256m -jar $DIR/iris-server-*.jar"

# SELinux requires runuser instead of su
if [ -x /sbin/runuser ]
then
	SU=/sbin/runuser
else
	SU=su
fi

function start() {
	pid=`pidof -s $PROG`
	if [ $pid ]
	then
		echo $"$SERVICE already running."
	else
		echo -n $"Starting $SERVICE:"
		$SU -l $USER -s /bin/sh -c "cd $DIR && $PROG $OPTIONS &" && success || failure
		RETVAL=$?
		[ "$RETVAL" = 0 ] && touch /var/lock/subsys/$SERVICE
		echo
	fi
}

function stop() {
	echo -n $"Stopping $SERVICE:"
	killproc $PROG -TERM
	RETVAL=$?
	[ "$RETVAL" = 0 ] && rm -f /var/lock/subsys/$SERVICE
	echo
}

case "$1" in
	start)
		start
		;;
	stop)
		stop
		;;
	restart)
		stop
		sleep 1
		start
		;;
	status)
		status $PROG
		RETVAL=$?
		;;
	*)
		echo $"Usage: $0 {start|stop|restart|status}"
		RETVAL=1
esac

exit $RETVAL
