# bin/iris.service
[Unit]
Description=IRIS Advanced Traffic Management System (ATMS)
After=postgresql.service

[Service]
ExecStart=/usr/bin/java -server -Xmx512m -jar /usr/share/java/iris-server-@@VERSION@@/iris-server-@@VERSION@@.jar
Type=simple
Restart=always
User=tms
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target
