# bin/iris
#!/bin/bash
#
# Init file for IRIS server daemon
#
# chkconfig: 2345 97 02
# description: IRIS server daemon
#
# processname: iris

# source function library
. /etc/rc.d/init.d/functions

RETVAL=0
SERVICE=$(basename $0)
PROG=iris_service
DIR=/usr/share/java/iris-server
USER=tms
OPTIONS="/usr/bin/java -server -Xmx512m -jar $DIR/iris-server-*.jar"

# SELinux requires runuser instead of su
if [ -x /sbin/runuser ]
then
	SU=/sbin/runuser
else
	SU=su
fi

function start() {
	pid=`pidof -s $PROG`
	if [ $pid ]
	then
		echo $"$SERVICE already running."
	else
		check_initialized
		if [ "$?" != "0" ]
		then
			echo "Have you tried 'iris_ctl init'?"
			exit 1
		fi
		prepare_client
		echo -n $"Starting $SERVICE:"
		$SU $USER -s /bin/sh -c "$PROG $OPTIONS &" && success || failure
		RETVAL=$?
		[ "$RETVAL" = 0 ] && touch /var/lock/subsys/$SERVICE
		echo
	fi
}

function check_initialized() {
	if [ ! -r /etc/iris/iris-server.keystore ]
	then
		echo "/etc/iris/iris-server.keystore not found."
		return 1
	fi
	if [ ! -r /etc/iris/iris-client.keystore ]
	then
		echo "/etc/iris/iris-client.keystore not found."
		return 1
	fi
	return 0
}

function stop() {
	clean_client
	echo -n $"Stopping $SERVICE:"
	killproc $PROG -TERM
	RETVAL=$?
	[ "$RETVAL" = 0 ] && rm -f /var/lock/subsys/$SERVICE
	echo
}

function prepare_client() {
	install -o apache -g apache -m 0444 /etc/iris/iris-client.keystore /var/www/html/iris-client/
	install -o apache -g apache -m 0444 /etc/iris/iris-client.properties /var/www/html/iris-client/
}

function clean_client() {
	rm -f /var/www/html/iris-client/iris-client.properties
	rm -f /var/www/html/iris-client/iris-client.keystore
}

case "$1" in
	start)
		start
		;;
	stop)
		stop
		;;
	restart)
		stop
		sleep 1
		start
		;;
	status)
		status $PROG
		RETVAL=$?
		;;
	*)
		echo $"Usage: $0 {start|stop|restart|status}"
		RETVAL=1
esac

exit $RETVAL
