# bin/honeybee/honeybee.service
[Unit]
Description=Honeybee service for IRIS
StartLimitIntervalSec=0

[Service]
ExecStart=/usr/local/bin/honeybee
Restart=always
RestartSec=2
User=tms

[Install]
WantedBy=multi-user.target
